# web-app/initscripts/lpc-ldap-server.conf
# Upstart configuration to manage the lpc-ldap-server service.

description "LPC LDAP server UI provider"

start on runlevel [2345]
stop on runlevel [016]

respawn

env APPDIR=/usr/lib/nodejs/lpc-ldap-server

pre-start script
    test -f $APPDIR/index.js || { stop; exit 0; }
end script

exec nodejs $APPDIR/index.js
